# system-config/receipt-printer-mqtt.service
[Unit]
Description=Quote Receipt Printer (MQTT)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/quotes
ExecStart=/usr/bin/python3 /home/pi/quotes/src/mqtt_print_subscriber.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
