# system-config/receipt-printer.service
[Unit]
Description=Quote Receipt Printer
After=network-online.target avahi-daemon.service
Wants=network-online.target

[Service]
Type=simple
# User configured for monalisa
User=monalisa
# Working directory must be absolute
WorkingDirectory=/home/monalisa/quotes
# Python path and app location
ExecStart=/usr/bin/python3 /home/monalisa/quotes/src/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
